# systemd/yodo_git-pull.service
[Unit]
Description=Git Pull Service

[Service]
Type=simple
WorkingDirectory=XXX
ExecStart=/usr/bin/git pull

