# systemd/yodo_git-pull.timer
[Unit]
Description=Run Git Pull Daily at 24:00

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
